# com.rosterfocus.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rosterfocus.agent</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Point this at the VENV python that has pyobjc-framework-EventKit
             installed (e.g. /path/to/roster-focus/.venv/bin/python). Do NOT use
             /usr/bin/python3 — the system Python lacks the bindings. The
             install.sh script fills these two paths in for you automatically. -->
        <string>/path/to/roster-focus/.venv/bin/python</string>
        <string>/path/to/roster-focus/rosterfocus.py</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/rosterfocus.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/rosterfocus.err.log</string>
</dict>
</plist>
